# deploy/core/supervisor/conf.d/crond.conf
[program:crond]
process_name=%(program_name)s
command=cron -f
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout